# mac/launchd/plists/com.eoin.backup-claude-memory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eoin.backup-claude-memory</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eoin/.local/bin/backup-claude-memory.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/eoin/.claude/projects</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/eoin/.local/bin/backup-claude-memory-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/eoin/.local/bin/backup-claude-memory-launchd.log</string>
</dict>
</plist>
